# files/etc/uci-defaults/voucher
#!/bin/sh

cat etc/nodogsplash.conf.captive-portal >> /etc/nodgsplash.conf
